# make_debs/install_scripts/adtn/postinstall

idplat=`grep -oP "id \= \K.+" $INIPATH`			#get default id from ini	
ipplat=`grep -oP "ip \=\s+\K\d*.\d*.\d*.\d*" $INIPATH`	#get default ip from ini
portplat=`grep -oP "port \=\s+\K\d*" $INIPATH` 		#get default port from ini

#Create system user for adtn
useradd -M -r -s /sbin/nologin $USER

sugIdPlat=`hostname`
sugIpPlat=`/sbin/ip r | grep default | cut -d \  -f 5 | xargs ip a l dev | grep inet\ | cut -d \  -f 6 | cut -d/ -f 1`
sugPortPlat=$portplat

echo -e "\e[92mStarting adtn configuration process..."
read -p "Set the adtn identifier (default is ${sugIdPlat}): " id
read -p "Set the ip (default is ${sugIpPlat}): " ip
read -p "Set the adtn port (default is ${sugPortPlat}): " port

if [ -n "$id" ]; then
    sed -i "s/id = $idplat/id = $id/g" $INIPATH
else
    sed -i "s/id = $idplat/id = $sugIdPlat/g" $INIPATH
fi

if [ -n "$ip" ]; then
    sed -i "s/ip = $ipplat/ip = $ip/g" $INIPATH
else
    sed -i "s/ip = $ipplat/ip = $sugIpPlat/g" $INIPATH
fi

if [ -n "$port" ]; then
  sed -i "s/port = $portplat/port = $port/g" $INIPATH
else
  sed -i "s/port = $portplat/port = $sugPortPlat/g" $INIPATH
fi

#Config adtn.ini permissions
chown $USER:$USER $INIPATH
chmod 644 $INIPATH
#Create data directory
mkdir -p $DATAPATH
chown -R $USER:$USER $DATAPATH
chmod 755 $DATAPATH 
#Create pid directory for adtn and set right permissions
mkdir -p $PIDPATH 
chown -R $USER:$USER $PIDPATH
chmod 755 $PIDPATH
#Make init.d script executable
chmod +x /etc/init.d/$MAINSRVC
#Add init.d script at startup
update-rc.d $MAINSRVC defaults > /dev/null 2>&1
#Add libraries to ldconfig
ldconfig
echo -e "Configuration process completed!\e[39m"
echo "Remember that you will be able to edit those"
echo "params and others in: ${INIPATH}"
exit 0;
